# deploy/systemd/local-rag-assistant-backup.timer
[Unit]
Description=Daily backup of local-rag-assistant vector database

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
